# deploy/hotspotai-scraper.service
[Unit]
Description=HotSpotAI Scraper Service
After=network.target hotspotai-backend.service

[Service]
Type=simple
User=hotspotai
Group=hotspotai
WorkingDirectory=/opt/hotspotai/HotSpotAI
Environment="PATH=/opt/hotspotai/HotSpotAI/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/hotspotai/HotSpotAI/.venv/bin/python -c "from core.scheduler import start_scheduler; start_scheduler()"
Restart=always
RestartSec=30

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hotspotai-scraper

# 安全
NoNewPrivileges=true
PrivateTmp=true

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
